# systemd/resource-kuma-collect.service
[Unit]
Description=resource-kuma collector (one shot)
After=network.target

[Service]
Type=oneshot
ExecStart=/opt/resource-kuma/collect.sh
User=root
